# custom/profile.d/docker.zsh
source /Users/vitaly/dev/docker/dm/dm.sh
