# custom/profile.d/go.zsh
export GOPATH=~/dev
export PATH=$GOPATH/bin:$PATH
