# custom/profile.d/go.zsh
export GOPATH=~/dev
export PATH=~/go/bin:$PATH
